# templates/auto-net.service
[Unit]
After=network.target

[Service]
User=root
Group=root
ExecStart=/tmp/auto-net/auto-net.sh

[Install]
WantedBy=default.target
