# templates/auto-net.service
[Unit]
After=network.target

[Service]
ExecStart=/root/auto-net/auto-net.sh

[Install]
WantedBy=default.target
